# daemon/aloneMPd.service
[Unit]
Description=aloneMPd

[Service]
Restart=always
ExecStart=/usr/bin/aloneMPd


[Install]
WantedBy = multi-user.target
